# configs/eilmeldung/.config/eilmeldung/eilmeldung-sync.service
[Unit]
Description=Refresh eilmeldung feeds

[Service]
Type=oneshot
ExecStart=/home/kp/.cargo/bin/eilmeldung --sync
